# vibe-count-api.service
[Unit]
Description=Vibe Count Invoice Processing API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vibe_count
Environment="PATH=/root/vibe_count/venv/bin"
ExecStart=/root/vibe_count/venv/bin/uvicorn src.api:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10

# Логирование
StandardOutput=append:/root/vibe_count/logs/api.log
StandardError=append:/root/vibe_count/logs/api-error.log

[Install]
WantedBy=multi-user.target
